# bwv_sync/bwv-sync.cron
35 5 * * * set -a; . /root/env.env; set +a;  /usr/local/bin/looplijsten_bwv_sync.sh
20 20 * * * set -a; . /root/env.env; set +a;  /usr/local/bin/looplijsten_bwv_sync.sh
